# deploy/systemd/hostvim-engine.service
[Unit]
Description=Hostvim Engine (hosting API)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Environment=HOSTVIM_HOME=__HOSTVIM_HOME__
Environment=HOSTVIM_CONFIG_DIR=/etc/hostvim
Environment=PANELSAR_HOME=__HOSTVIM_HOME__
Environment=PANELSAR_CONFIG_DIR=/etc/hostvim
ExecStart=__ENGINE_BINARY__
Restart=on-failure
RestartSec=5
# Güvenlik: yalnızca yerel panel erişsin (nginx üzerinden değil; doğrudan loopback)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
