# deploy/systemd/panelsar-engine.service
[Unit]
Description=Panelsar Engine (hosting API)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
Environment=PANELSAR_HOME=__PANELSAR_HOME__
Environment=PANELSAR_CONFIG_DIR=/etc/panelsar
ExecStart=__ENGINE_BINARY__
Restart=on-failure
RestartSec=5
# Güvenlik: yalnızca yerel panel erişsin (nginx üzerinden değil; doğrudan loopback)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
